# systemd/scimitar-diag-hid@.service
[Unit]
Description=Corsair Scimitar raw HID report capture (Layer B) on /dev/%i
Documentation=https://github.com/jmylchreest/corsair-investigate
# Instances are started by udev (SYSTEMD_WANTS) when a Corsair hidraw
# node appears, and stop automatically when it vanishes on unplug.
BindsTo=dev-%i.device
After=dev-%i.device

[Service]
Type=simple
User=root
EnvironmentFile=-/etc/scimitar-diag/scimitar-diag.env
ExecStart=/usr/local/bin/scimitar_hid -d /dev/%i -o ${HID_LOG_DIR}/hid_reports.%i.log --rotate-mb ${ROTATE_MB}
Restart=on-failure
RestartSec=5

ProtectSystem=strict
ReadWritePaths=/var/log/scimitar-diag
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
